# systemd/bot-brain.service
[Unit]
Description=Telegram Bot Army — Brain Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/projekte/telegram-bot-army
EnvironmentFile=/root/projekte/telegram-bot-army/.env
Environment=PATH=/root/.nvm/versions/node/v24.16.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/python3 /root/projekte/telegram-bot-army/bots/brain.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
